# com.shawn.email-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shawn.email-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/shawnshlee/1_Claude_Code/email2ppt/venv/bin/python</string>
        <string>/Users/shawnshlee/1_Claude_Code/email2ppt/watcher.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/shawnshlee/1_Claude_Code/email2ppt</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SHOMERY_DISABLE_TELEGRAM</key>
        <string>1</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>StandardOutPath</key>
    <string>/Users/shawnshlee/1_Claude_Code/email2ppt/watcher.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shawnshlee/1_Claude_Code/email2ppt/watcher.stderr.log</string>
</dict>
</plist>
