# com.shawn.email-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shawn.email-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/shawnlee/telegram-bridge/venv/bin/python</string>
        <string>/Users/shawnlee/telegram-bridge/watcher.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/shawnlee/telegram-bridge</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>/Users/shawnlee/telegram-bridge/watcher.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shawnlee/telegram-bridge/watcher.stderr.log</string>
</dict>
</plist>
