# com.shawn.email-retention.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shawn.email-retention</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/shawnlee/telegram-bridge/venv/bin/python</string>
        <string>/Users/shawnlee/telegram-bridge/retention_sweep.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/shawnlee/telegram-bridge</string>
    <key>RunAtLoad</key>
    <false/>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/shawnlee/telegram-bridge/retention_sweep.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shawnlee/telegram-bridge/retention_sweep.stderr.log</string>
</dict>
</plist>
